# weaveworld_selenium/fs/etc/supervisor/conf.d/selenium.conf
[program:selenium]
command=/usr/share/graalvm/bin/java -jar /opt/selenium/selenium-server.jar standalone
environment=DISPLAY=":0",HOME=%(ENV_HOME)s,USER=%(ENV_USER)s
autostart=true
autorestart=false